# rfid-access-control.service
[Unit]
Description=RFID Access Control System
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/rfid-access-control
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PYTHONPATH=/home/pi/rfid-access-control
ExecStart=/usr/bin/python3 /home/pi/rfid-access-control/integrated_access_camera.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/rfid-access-control

[Install]
WantedBy=multi-user.target
